# yggdrasilgate/supervisord.conf
[supervisord]
logfile=/dev/null
nodaemon=true
user=root

[program:yggdrasil]
command=yggdrasil -useconffile ./yggdrasil.conf
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:app-forward]
command=socat TCP6-LISTEN:80,fork,forever,reuseaddr TCP4:app:80